# DEBIAN/postinst
#!/bin/sh

set -e

# generated by other debhelper scripts.
service="
[Unit]
Description = monacoind deamon

[Service]
ExecStart = monacoind -deamon
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
"
echo "Finished installing monacoind."
echo "Adding monacoind to service.."
sudo echo "$service" > /etc/systemd/system/monacoind.service
sudo systemctl enable monacoind
echo "Added monacoind to service! execute  sudo systemctl monacoind start  to start!"
echo "and make sure you wrote config file before starting monacoind! or else, monacoind will give you an error!"
echo "Now, this installation is over! Have fan!"

echo "Copyright for Monacoind:"
echo "Copyright (C) 2009-2018 The Bitcoin Core developers
Copyright (C) 2011-2018 The Litecoin Core developers
Copyright (C) 2011-2018 Dr. Kimoto Chan
Copyright (C) 2011-2018 The DigiByte developers
Copyright (C) 2014-2018 The Dash developers
Copyright (C) 2014-2018 The Vertcoin Developers
Copyright (C) 2013-2019 The Monacoin Core developers"
echo "Copyright for this installer:"
echo "Copyright (C) 2018 raspi0124"
echo "This installer is licenced with GPLv3."

exit 0
